# backend/launchd/com.theo.othello-v2-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.theo.othello-v2-worker</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/theo/signal-intelligence/othello_v2/backend/.venv/bin/python</string>
    <string>worker.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/theo/signal-intelligence/othello_v2/backend</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>OTHELLO_INTERNAL_SCHEDULER</key>
    <string>true</string>
    <key>OTHELLO_WORKER_BOOTSTRAP_MODE</key>
    <string>none</string>
    <key>OTHELLO_WORKER_ENABLE_INGESTION</key>
    <string>true</string>
    <key>OTHELLO_WORKER_ENABLE_TRANSLATIONS</key>
    <string>false</string>
    <key>OTHELLO_WORKER_ENABLE_ANALYTICS</key>
    <string>false</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string>/tmp/othello_v2_worker.launchd.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/othello_v2_worker.launchd.log</string>
</dict>
</plist>
